# var/config-static/etc/cron.d/clamav-freshclam
57 */24 * * *    clamav [ -x /usr/bin/freshclam ] && /usr/bin/freshclam >/dev/null
