# var/config-static/var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.IebCjY/crontab installed on Thu Nov  2 13:45:07 2006)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
0 1 1 * * /usr/sbin/linuxmuster-backup --full
0 2 9,17,25 * * /usr/sbin/linuxmuster-backup --diff
0 3 2-8,10-16,18-24,26-31 * * /usr/sbin/linuxmuster-backup --inc
